# deploy/cron/ats-auto-login.cron
# /etc/cron.d/ats-auto-login
# Daily Zerodha auto-login. Drops to /etc/cron.d/ via the install script.
#
# Schedule: 03:20 UTC = 08:50 IST, Mon-Fri (Indian markets are closed weekends).
# The script itself adds 0-60s jitter to avoid trivial pattern detection.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

20 3 * * 1-5 root /usr/local/bin/ats-morning-check.sh >> /var/log/ats-morning-check.log 2>&1
